# services/memcache.conf
[program:memcache]
command = /usr/bin/memcached -p11211 -l0.0.0.0 -uroot -m%(ENV_MEMCACHE_MAX_MEMORY)s -I%(ENV_MEMCACHE_SLAB_SIZE)s
priority = 999
redirect_stderr=true
umask=007
